# ENTER/pkgs/conda-25.7.0-py313hca03da5_0/etc/profile.d/conda.sh
version https://git-lfs.github.com/spec/v1
oid sha256:b5c510d50128e10a3a40e8ba941d028f2ff9fb7580d7ceac5a2439d6baef9ad3
size 3379
