# ENTER/pkgs/libmamba-2.0.5-h15e39b3_1/etc/profile.d/mamba.sh
version https://git-lfs.github.com/spec/v1
oid sha256:2978300ce8266491a27675166f3fadb53c1d26b5f212bbb0bf33993012c623f8
size 1008
