# ENTER/etc/profile.d/conda.sh
version https://git-lfs.github.com/spec/v1
oid sha256:d58abf9f629b9b9e524cede9acda254da9ea7e715ccd2b8474fe0f924da88643
size 2571
