# ENTER/etc/profile.d/mamba.sh
version https://git-lfs.github.com/spec/v1
oid sha256:8fc2c6ae53e88430c66f665b2e08d8f3512659e89792e3106779828de2e508e8
size 402
